# kubeedge-configuration/systemd service/edgeIP.service
[Unit]
Description=Manage public IP of edge node.

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/edgeIPcheck.sh