# kubeedge-configuration/systemd service/edgeIP.timer
[Unit]
Description=Edge node public IP service timer

[Timer]
OnUnitActiveSec=12h
OnBootSec=1h

[Install]
WantedBy=timers.target